# nextcloud-rag-ui.service
[Unit]
Description=Nextcloud RAG web UI
After=network-online.target llama-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/nextcloud-rag
ExecStart=/opt/nextcloud-rag/.venv/bin/python /opt/nextcloud-rag/rag_app.py serve --host 0.0.0.0 --port 7860
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
